# MM_13/page.service
[Unit]
Description=My Java WebAPP

[Service]
User=root
WorkingDirectory=/kursova
ExecStart=/kursova/start

[Install]
WantedBy=multi-user.target
